# crontab

00 22 * * * root /app/update_mongo.sh >> /var/log/cron.log 2>&1
